# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

; ── Caddy (reverse proxy + static landing page) ───────────────────────────────
[program:caddy]
command=/usr/bin/caddy run --config /config/caddy/Caddyfile --adapter caddyfile
user=root
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/var/log/caddy.log
stderr_logfile=/var/log/caddy.err
environment=HOME="/root",XDG_DATA_HOME="/config/caddy/.local"

; ── JupyterLab ────────────────────────────────────────────────────────────────
[program:jupyterlab]
command=/home/engineering/.venv/bin/jupyter lab
    --config=/home/engineering/.jupyter/jupyter_lab_config.py
    --no-browser
user=engineering
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/jupyterlab.log
stderr_logfile=/var/log/jupyterlab.err
environment=
    HOME="/home/engineering",
    PATH="/home/engineering/.venv/bin:/usr/local/bin:/usr/bin:/bin",
    VIRTUAL_ENV="/home/engineering/.venv"

; ── code-server (VS Code) ─────────────────────────────────────────────────────
[program:code-server]
command=/usr/bin/code-server
    --bind-addr 127.0.0.1:8081
    --auth none
    --user-data-dir /config/code-server/user-data
    --extensions-dir /config/code-server/extensions
    /home/engineering/work
user=engineering
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/code-server.log
stderr_logfile=/var/log/code-server.err
environment=
    HOME="/home/engineering",
    PATH="/home/engineering/.venv/bin:/usr/local/bin:/usr/bin:/bin"
# [program:code-server]
# command=/usr/bin/code-server
#     --bind-addr 127.0.0.1:8081
#     --auth none
#     --user-data-dir /config/code-server/user-data
#     --extensions-dir /config/code-server/extensions
#     --app-name "VS Code"
#     --base-path /code
#     /home/engineering/work
# user=engineering
# autostart=true
# autorestart=true
# startsecs=5
# startretries=5
# stdout_logfile=/var/log/code-server.log
# stderr_logfile=/var/log/code-server.err
# environment=
#     HOME="/home/engineering",
#     PATH="/home/engineering/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# ; ── Tabula ────────────────────────────────────────────────────────────────────
# [program:tabula]
# command=/usr/bin/java -Dfile.encoding=utf-8 -Xms256M -Xmx1024M
#     -Dtabula.openBrowser=false
#     -Dtabula.port=9999
#     -jar /usr/local/bin/tabula.jar
# user=engineering
# autostart=true
# autorestart=true
# startsecs=3
# startretries=5
# stdout_logfile=/var/log/tabula.log
# stderr_logfile=/var/log/tabula.err
# environment=HOME="/home/engineering"

; ── ttyd (web terminal) ───────────────────────────────────────────────────────
[program:ttyd]
command=/usr/local/bin/ttyd
    --port 7681
    --interface 127.0.0.1
    --writable
    login -f engineering
user=root
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/var/log/ttyd.log
stderr_logfile=/var/log/ttyd.err
environment=HOME="/home/engineering",TERM="xterm-256color"

; ── Tabula web app ────────────────────────────────────────────────────────────
[program:tabula]
command=/usr/bin/java -Dfile.encoding=utf-8 -Xms256M -Xmx1024M -Dwarbler.port=9999 -jar /usr/local/bin/tabula-web/tabula/tabula.jar
user=engineering
autostart=true
autorestart=true
startsecs=8
startretries=5
stdout_logfile=/var/log/tabula.log
stderr_logfile=/var/log/tabula.err
environment=HOME="/home/engineering"